# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/local/bin/service.sh

[program:borg]
command=/usr/bin/borg serve --restrict-to-path /home/backup/repository
